# uap3-player.service
[Unit]
Description=UAP Dog Whistle Player
After=bluetooth.target sound.target
Wants=bluetooth.target

[Service]
Type=simple
User=pi
ExecStart=/usr/bin/python3 /home/pi/uap3-player.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
